# install-scripts/systemd/revelnest-yggdrasil.service
[Unit]
Description=RevelNest Yggdrasil Mesh Network
After=network.target
Wants=network.target
Documentation=https://revelnest.chat/docs

[Service]
Type=simple
# El ejecutable se instalará en /usr/lib/revelnest/yggdrasil
ExecStart=/usr/lib/revelnest/yggdrasil -useconffile /etc/revelnest/yggdrasil.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Usuario y grupo dedicados para seguridad
User=yggdrasil
Group=yggdrasil
# Directorio de trabajo donde se almacenan datos persistentes
WorkingDirectory=/var/lib/revelnest
# Capabilities necesarias para crear interfaces TUN sin root completo
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/revelnest
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectHostname=true

# Configuración de logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=revelnest-yggdrasil

# Límites de recursos
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target